# supervisord.conf
[supervisord]
nodaemon=true

[program:npm]
command=npm run start
autostart=true
autorestart=true
stderr_logfile=/var/log/npm.err.log
stdout_logfile=/var/log/npm.out.log
